# roles/voctomix/files/systemd/videoteam-voctocore.service
[Unit]
Description=Voctomix core Daemon

[Service]
ExecStart=/usr/bin/voctocore -v
Restart=always
RestartSec=1s

[Install]
WantedBy=default.target
