# roles/voctomix/files/systemd/videoteam-capture-blackmagic-hdmi@.service
[Unit]
Description=Voctomix source HDMI blackmagic from card %I
After=videoteam-voctocore.service
Requires=videoteam-voctocore.service
StartLimitIntervalSec=0
Conflicts=xlogout.target

[Service]
ExecStart=/usr/bin/voctomix-ingest -v --video-source blackmagic --video-attribs 'device-number=%I connection=2' --port 1000%I
Restart=always
RestartSec=1s

[Install]
WantedBy=xlogin.target
