# roles/opsis/files/videoteam-ingest.path
[Path]
PathExists=/dev/hdmi2usb/by-num/all0/video

[Install]
WantedBy=paths.target
