# roles/xorg/files/videoteam-konversation.service
[Unit]
Description=Konversation IRC Client
Conflicts=xlogout.target

[Service]
ExecStart=/usr/bin/konversation

[Install]
WantedBy=xlogin.target
