# roles/xorg/files/xlogin.target
[Unit]
Description=X Login (via lightdm)
